# X/.config/systemd/user/urxvtd.service
[Unit]
Description=URvxt daemon service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/urxvtd 
Restart=always
RestartSec=3

