# feh-background/.config/systemd/user/feh-background.timer
[Unit]
Description=Change walpaper each 2 minutes

[Timer]
OnUnitActiveSec=2m

[Install]
WantedBy=default.target

