# feh-background/.config/systemd/user/feh-background@.timer
[Unit]
Description=Change walpaper each %I

[Timer]
OnUnitActiveSec=%i
Unit=feh-background.service

[Install]
WantedBy=default.target

